# systemd/wallsplash.service
[Unit]
Description=Change wallpaper from unsplash

[Service]
ExecStart=%h/unspla.sh

[Install]
WantedBy=default.target
